# php/s6/services.d/php-fpm/run
#!/command/with-contenv sh
set -eu
exec 2>&1

exec docker-php-entrypoint php-fpm -F
